# Systemd_Units/diskstatus.service
[Unit]
Description=Disk status checker

[Service]
User=root
Group=root
WorkingDirectory=/root
Type=simple
Type=oneshot
ExecStart=/root/diskstatus.sh
